# barman/crontab
*/1 * * * * * /usr/bin/barman cron --keep-descriptors > /proc/1/fd/1 2> /proc/1/fd/2
